# systemd/trigger-monitor.service
[Unit]
Description=Superalignment Trigger Monitor (One-shot)
After=network.target

[Service]
Type=oneshot
User=lizthedeveloper_gmail_com
WorkingDirectory=/home/lizthedeveloper_gmail_com/src/superalignment-chatroom
ExecStart=/usr/bin/python3 /home/lizthedeveloper_gmail_com/src/superalignment-chatroom/scripts/trigger-monitor.py --once
StandardOutput=append:/home/lizthedeveloper_gmail_com/src/superalignment-chatroom/logs/trigger-monitor-service.log
StandardError=append:/home/lizthedeveloper_gmail_com/src/superalignment-chatroom/logs/trigger-monitor-service.log

# Environment
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Security
PrivateTmp=yes
NoNewPrivileges=yes

# Restart on failure
Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
