# systemd/trigger-monitor.timer
[Unit]
Description=Run Superalignment Trigger Monitor every 5 minutes
Requires=trigger-monitor.service

[Timer]
# Run every 5 minutes after the service completes
OnUnitActiveSec=5min
# Start 1 minute after boot
OnBootSec=1min
# Run immediately if timer was missed
Persistent=yes

[Install]
WantedBy=timers.target
